# docker/supervisor-config/node.sv.conf
[program:node]
command = /usr/bin/node /var/nodejs/app.js
